# route-optimizer.service
[Unit]
Description=Innofoods Route Optimizer Web
After=network.target

[Service]
Type=simple
User=sam
WorkingDirectory=/home/sam/route-optimizer
ExecStart=/usr/bin/python3 -u /home/sam/route-optimizer/route_optimizer_web.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
